# nginx-conf/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /usr/src/cryptex-api/nginx-conf/uwsgi.ini

[program:nginx-app]
command = /usr/sbin/nginx